# docker/test/supervisord/supervisor.d/swoole.ini
; 是被管理的进程配置参数，“swoole”是进程的名称
[program:swoole]
; 程序启动命令
command=php think swoole
; 目录
directory=/app
; 在supervisord启动的时候也自动启动
autostart=true
; 程序退出后自动重启，可选值：[unexpected,true,false]，默认为unexpected，表示进程意外杀死后才重启
autorestart=true
; stdout 日志文件，需要注意当指定目录不存在时无法正常启动，所以需要手动创建目录（supervisord 会自动创建日志文件）
stdout_logfile=/dev/stdout
; stdout 日志文件大小，默认50MB
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0